# supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/stdout
logfile_maxbytes=0
loglevel=info

; Virtual framebuffer — Wine needs a display even for a headless server
[program:xvfb]
command=Xvfb :99 -screen 0 1024x768x16 -ac -nolisten tcp
autostart=true
autorestart=true
priority=1
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; VNC server — streams the virtual display so you can see DServer's window
[program:x11vnc]
command=bash -c 'sleep 3; exec x11vnc -display :99 -nopw -forever -shared -rfbport 5900 -noxdamage -quiet'
autostart=true
autorestart=true
priority=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; noVNC websocket proxy — serves the browser-based VNC viewer on port 6080
[program:novnc]
command=websockify --web /usr/share/novnc/ 0.0.0.0:6080 localhost:5900
autostart=true
autorestart=true
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; Flask web management interface
[program:web]
command=python3 /web/app.py
directory=/web
autostart=true
autorestart=true
priority=20
environment=DISPLAY=":99",HOME="/root",WINEPREFIX="/opt/wineprefix"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
